# systemd/freezer-web.service
[Unit]
Description=Freezer Monitor — Web Dashboard Service
After=network.target
RequiresMountsFor=/data

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 -u /opt/freezerpi/web_server.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
ReadWritePaths=/run/freezerpi /run/freezer_db /data
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
