# systemd/icebox-web.service
[Unit]
Description=IceboxHero — Web Dashboard Service
After=network.target
RequiresMountsFor=/data

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 -u /opt/iceboxhero/web_server.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
ReadWritePaths=/run/iceboxhero /run/icebox_db /data
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
